# systemd/pjtaudi-watchdog.timer
[Unit]
Description=Run PJTAudiBot watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=pjtaudi-watchdog.service

[Install]
WantedBy=timers.target
